# package/pifii/selectchannel/files/selectchannel.init
#!/bin/sh /etc/rc.common

START=99

local pro="/usr/sbin/selectchannel"

start(){
    test -f "$pro" && lua "$pro" &
}
